# deploy/jester.service
# systemd unit for hermes-jester on an exe.dev VM.
# Install to /etc/systemd/system/jester.service (or ~/.config/systemd/user/ for a
# user service), then: systemctl enable --now jester
#
# Adjust User, WorkingDirectory, and the EnvironmentFile path to match your VM.

[Unit]
Description=hermes-jester queue service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/hermes-jester
EnvironmentFile=/home/ubuntu/hermes-jester/.env
ExecStart=/home/ubuntu/hermes-jester/.venv/bin/uvicorn jester.app:create_app --factory --host 0.0.0.0 --port 8077
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
